# conf/htb-presence.service
[Unit]
Description=HTB Presence

[Service]
ExecStart=/usr/local/bin/htb-presence.py
User=

[Install]
WantedBy=default.target


